# sdbootutil-update-predictions.service
[Unit]
Description=Update TPM predictions
ConditionSecurity=tpm2

[Service]
Type=oneshot
RemainAfterExit=yes
KeyringMode=shared
PrivateTmp=yes
ExecStartPre=/usr/bin/sh -c 'mkdir "/tmp/$INVOCATION_ID"; cp -a /var/lib/pcrlock.d /tmp/$INVOCATION_ID/pcrlock.d'
ExecStart=/usr/bin/sdbootutil -v update-predictions
ExecStopPost=/usr/bin/sh -c '\
    if [ "$SERVICE_RESULT" != "success" ]; then \
        echo "Command failed. Reason: $SERVICE_RESULT (Status: $EXIT_STATUS). Restoring backup..."; \
        rm -rf /var/lib/pcrlock.d/*; \
        cp -a "/tmp/$INVOCATION_ID/pcrlock.d/." /var/lib/pcrlock.d/; \
    fi; \
    rm -rf "/tmp/$INVOCATION_ID/pcrlock.d"'
ImportCredential=sdbootutil-update-predictions.*

[Install]
WantedBy=default.target
