# sdbootutil-update-predictions.service
[Unit]
Description=Update TPM predictions
ConditionSecurity=tpm2

[Service]
Type=oneshot
ExecStart=/usr/bin/sdbootutil update-predictions
